# com.policeScraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.policeScraper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/melkor/Documents/Projects/policeScrapper/run_scraper.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>WorkingDirectory</key>
    <string>/Users/melkor/Documents/Projects/policeScrapper</string>
    <key>StandardOutPath</key>
    <string>/Users/melkor/Documents/Projects/policeScrapper/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/melkor/Documents/Projects/policeScrapper/logs/launchd.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist> 